# install/roles/prometheus/templates/prometheus.conf.j2
# Run prometheus

start on startup

chdir {{ prometheus_dist_dir }}/prometheus-{{prometheus_version}}.linux-amd64

script
   ./prometheus -config.file /etc/prometheus/prometheus.yml
end script

